# systemd.trump-truth-bot.service
[Unit]
Description=Trump Truth Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/ubuntu/.openclaw/workspace/trump-truth-bot
ExecStart=/home/ubuntu/nvm/versions/node/v24.14.0/bin/node /home/ubuntu/.openclaw/workspace/trump-truth-bot/src/index.js
Restart=always
RestartSec=5
User=ubuntu
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
